# os/rootfs/etc/init.d/wraith-samba
#!/sbin/openrc-run
# wraith-samba - Restore network mounts (CIFS + NFS) from config disk
#
# Reads mount definitions from /wraith/config/samba.json and mounts
# each share. Mounts go to /mnt/<name> so they can be
# bind-mounted into Docker containers.

description="Restore WraithOS network mounts"

# Load configuration
. /etc/conf.d/wraith

depend() {
    need wraith-disks wraith-network
    before wraith-ui
}

# Parse simple JSON arrays using sed/grep.
# This handles the samba.json format:
# {"mounts":[{"mountpoint":"/mnt/x","server":"y","share":"z","username":"u","password":"p","options":"o"}]}
# Also supports legacy format with "name" field instead of "mountpoint".
#
# For each mount entry, extracts fields and calls mount_share.
parse_and_mount() {
    local config_file="$1"

    [ -f "$config_file" ] || return 0

    # Extract individual mount objects (one per line)
    # This uses a simple approach: split on },{ to get individual entries
    local mounts_json
    mounts_json=$(sed 's/.*"mounts":\[//;s/\].*//' "$config_file" \
        | sed 's/},{/}\n{/g')

    [ -z "$mounts_json" ] && return 0

    echo "$mounts_json" | while IFS= read -r entry; do
        [ -z "$entry" ] && continue

        local name server share username password mount_opts mount_type
        local mountpoint_val

        # Extract mountpoint (new format) or name (legacy format)
        mountpoint_val=$(echo "$entry" | sed -n 's/.*"mountpoint"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p')
        if [ -n "$mountpoint_val" ]; then
            name=$(basename "$mountpoint_val")
        else
            # Legacy fallback: use "name" field directly
            name=$(echo "$entry" | sed -n 's/.*"name"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p')
        fi
        server=$(echo "$entry"   | sed -n 's/.*"server"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p')
        share=$(echo "$entry"    | sed -n 's/.*"share"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p')
        username=$(echo "$entry" | sed -n 's/.*"username"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p')
        password=$(echo "$entry" | sed -n 's/.*"password"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p')
        mount_opts=$(echo "$entry" | sed -n 's/.*"options"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p')
        mount_type=$(echo "$entry" | sed -n 's/.*"type"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p')

        # Default to cifs if type not specified (backward compat)
        [ -z "$mount_type" ] && mount_type="cifs"

        if [ -z "$name" ] || [ -z "$server" ] || [ -z "$share" ]; then
            ewarn "Skipping incomplete mount entry"
            continue
        fi

        if [ "$mount_type" = "nfs" ]; then
            mount_nfs "$name" "$server" "$share" "$mount_opts"
        else
            mount_share "$name" "$server" "$share" "$username" "$password" "$mount_opts"
        fi
    done
}

# Mount a single CIFS share.
# Args: name server share [username] [password] [extra_options]
mount_share() {
    local name="$1"
    local server="$2"
    local share="$3"
    local username="$4"
    local password="$5"
    local extra_opts="$6"

    local mount_point="$WRAITH_SAMBA_MOUNT_DIR/$name"
    local unc="//${server}/${share}"

    # Skip if already mounted
    if mountpoint -q "$mount_point" 2>/dev/null; then
        einfo "Samba mount '$name' already mounted"
        return 0
    fi

    mkdir -p "$mount_point"

    # Build mount options
    local opts="noatime,noperm,iocharset=utf8,file_mode=0777,dir_mode=0777"
    if [ -n "$username" ]; then
        opts="${opts},username=${username}"
        if [ -n "$password" ]; then
            opts="${opts},password=${password}"
        fi
    else
        opts="${opts},guest"
    fi
    if [ -n "$extra_opts" ]; then
        opts="${opts},${extra_opts}"
    fi

    einfo "Mounting Samba share: $unc -> $mount_point"
    if mount -t cifs -o "$opts" "$unc" "$mount_point"; then
        einfo "  Mounted successfully"
    else
        ewarn "  Failed to mount $unc -- share may be unavailable"
        # Don't fail the service; other mounts may succeed
    fi
}

# Mount a single NFS share.
# Args: name server export [extra_options]
mount_nfs() {
    local name="$1"
    local server="$2"
    local export_path="$3"
    local extra_opts="$4"

    local mount_point="$WRAITH_SAMBA_MOUNT_DIR/$name"
    local source="${server}:${export_path}"

    # Skip if already mounted
    if mountpoint -q "$mount_point" 2>/dev/null; then
        einfo "NFS mount '$name' already mounted"
        return 0
    fi

    mkdir -p "$mount_point"

    local opts="noatime,nfsvers=4"
    if [ -n "$extra_opts" ]; then
        opts="${opts},${extra_opts}"
    fi

    einfo "Mounting NFS share: $source -> $mount_point"
    if mount -t nfs -o "$opts" "$source" "$mount_point"; then
        einfo "  Mounted successfully"
    else
        ewarn "  Failed to mount $source -- share may be unavailable"
    fi
}

start() {
    ebegin "Restoring network mounts"

    local config_file="$WRAITH_CONFIG_DIR/samba.json"

    if [ ! -f "$config_file" ]; then
        einfo "No network mount configuration found"
        eend 0
        return 0
    fi

    # Check for empty mounts array
    if grep -q '"mounts":\[\]' "$config_file" 2>/dev/null; then
        einfo "No network mounts configured"
        eend 0
        return 0
    fi

    mkdir -p "$WRAITH_SAMBA_MOUNT_DIR"
    parse_and_mount "$config_file"

    eend 0
}

stop() {
    ebegin "Unmounting network shares"

    # Unmount everything under the network mount directory
    if [ -d "$WRAITH_SAMBA_MOUNT_DIR" ]; then
        for mnt in "$WRAITH_SAMBA_MOUNT_DIR"/*/; do
            [ -d "$mnt" ] || continue
            if mountpoint -q "$mnt" 2>/dev/null; then
                umount -l "$mnt" 2>/dev/null || true
            fi
        done
    fi

    eend 0
}

status() {
    local count=0

    if [ -d "$WRAITH_SAMBA_MOUNT_DIR" ]; then
        for mnt in "$WRAITH_SAMBA_MOUNT_DIR"/*/; do
            [ -d "$mnt" ] || continue
            if mountpoint -q "$mnt" 2>/dev/null; then
                local name
                name=$(basename "$mnt")
                einfo "Samba mount '$name': mounted"
                count=$((count + 1))
            fi
        done
    fi

    einfo "Total active Samba mounts: $count"
    return 0
}
