# os/rootfs/etc/init.d/wraith-docker
#!/sbin/openrc-run
# wraith-docker - Deploy and manage the docker-compose stack for WraithOS
#
# The stock Alpine 'docker' service handles the actual dockerd lifecycle.
# daemon.json is written by wraith-disks (which runs before docker).
# This service deploys the compose stack after Docker is ready.

description="Deploy WraithOS docker-compose stack"

# Load configuration
. /etc/conf.d/wraith

depend() {
    need docker wraith-disks
    before wraith-samba
}

wait_for_docker() {
    local max_wait=30
    local waited=0
    while [ $waited -lt $max_wait ]; do
        if docker info >/dev/null 2>&1; then
            return 0
        fi
        sleep 1
        waited=$((waited + 1))
    done
    return 1
}

start() {
    ebegin "WraithOS compose stack"

    local compose_dir="$WRAITH_CONFIG_DIR/compose"
    local compose_file="$compose_dir/docker-compose.yml"

    if [ ! -f "$compose_file" ]; then
        einfo "No compose stack found -- skipping deploy"
        eend 0
        return 0
    fi

    einfo "Waiting for Docker daemon..."
    if ! wait_for_docker; then
        ewarn "Docker daemon not ready after 30s -- skipping compose deploy"
        eend 1
        return 1
    fi

    # Check if requireMounts is enabled in compose settings
    local settings_file="$compose_dir/settings.json"
    if [ -f "$settings_file" ]; then
        local require_mounts
        require_mounts=$(sed -n 's/.*"requireMounts"[[:space:]]*:[[:space:]]*\(true\|false\).*/\1/p' "$settings_file")
        if [ "$require_mounts" = "true" ]; then
            local samba_file="$WRAITH_CONFIG_DIR/samba.json"
            if [ -f "$samba_file" ]; then
                local all_mounted=true
                # Extract mount points from samba.json (flatten first — JSON is indented)
                for mp in $(tr -d '\n' < "$samba_file" | grep -o '"mountpoint"[[:space:]]*:[[:space:]]*"[^"]*"' | sed 's/.*"\([^"]*\)"$/\1/'); do
                    if ! grep -q " $mp " /proc/mounts 2>/dev/null; then
                        ewarn "Required mount not available: $mp"
                        all_mounted=false
                    fi
                done
                if [ "$all_mounted" = "false" ]; then
                    ewarn "Require mounts is enabled but not all mounts are available -- skipping compose deploy"
                    eend 1
                    return 1
                fi
            fi
        fi
    fi

    einfo "Deploying docker-compose stack from config disk"

    local env_flag=""
    if [ -f "$compose_dir/.env" ]; then
        env_flag="--env-file $compose_dir/.env"
    fi

    # shellcheck disable=SC2086
    docker compose -f "$compose_file" $env_flag up -d 2>&1 \
        | while read -r line; do einfo "  compose: $line"; done

    # Note: can't check PIPESTATUS on busybox sh; compose errors
    # are logged by the pipe above and visible in wraith-ui logs.

    einfo "Compose stack deployed successfully"
    eend 0
}

stop() {
    ebegin "Stopping WraithOS compose stack"

    local compose_file="$WRAITH_CONFIG_DIR/compose/docker-compose.yml"
    if [ -f "$compose_file" ] && docker info >/dev/null 2>&1; then
        einfo "Stopping compose stack..."
        docker compose -f "$compose_file" down --timeout 30 2>/dev/null || true
    fi

    eend 0
}

status() {
    if docker info >/dev/null 2>&1; then
        local containers
        containers=$(docker ps -q 2>/dev/null | wc -l)
        einfo "Docker: running ($containers containers active)"
        return 0
    else
        eerror "Docker: not running"
        return 1
    fi
}
